# systemd/lxc-auto-update.service
[Unit]
Description=Proxmox VE LXC Auto Update (running containers only)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/lxc-auto-update.sh
